# celery.supervisor.conf
; ==================================
;  celery worker supervisor example
; ==================================

[program:celery]
; Set full path to celery program if using virtualenv
command=/home/jqr/venv/masque/bin/celery -A tasks worker --loglevel=DEBUG

; Alternatively,
;command=celery --app=your_app.celery:app worker --loglevel=INFO -n worker.%%h
; Or run a script
;command=celery.sh

environment=
    PYTHONPATH="/home/jqr/masque-api-test",
    LOG_PATH="/home/jqr/masque-api-test/log",
    LOG_FILE_LEVEL="DEBUG",
    FLASK_DEBUG=0,
    MONGO_DB="masque",
    REDIS_DB="1",
    SMS_IKEY="",
    SMS_AKEY="",
    SMS_TEMPLATE_CODE="SMS_6770229",
    SMS_FREE_SIGN_NAME="身份验证",
    ALI_IKEY="",
    ALI_AKEY="",


directory=/home/jqr/masque-api-test
user=nobody
numprocs=1
stdout_logfile=/var/log/supervisor/celery_worker.log
redirect_stderr=true
autostart=true
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; When resorting to send SIGKILL to the program to terminate it
; send SIGKILL to its whole process group instead,
; taking care of its children as well.
killasgroup=true

; Set Celery priority higher than default (999)
; so, if rabbitmq is supervised, it will start first.
priority=1000
